# bash-scripts/mount-sd.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/mount-sd

[Install]
WantedBy=multi-user.target 
