# bash-scripts/wallpaper.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/wallpaper

[Install]
WantedBy=multi-user.target 
