# pmos-init-files/try
#!/sbin/openrc-run

name="try daemon"
command="/home/jonathan/Documents/phone/bash-scripts/try-daemon.sh"
pidfile="/home/jonathan/Documents/phone/bash-scripts/try-daemon.pid"
command_background=true

depend() {
	need net
}

