# pmos-init-files/Wallpaper
#!/sbin/openrc-run

name="Wallpaper daemon"
command="/home/jonathan/Documents/phone/bash-scripts/feh-change-background.sh"
pidfile="/home/jonathan/Documents/phone/bash-scripts/wallpaper.pid"
command_background=true

depend() {
	need net
}



